# packaging/systemd/deepiri-git-handshake.service
[Unit]
Description=Deepiri Git Handshake clone transport daemon
After=network.target

[Service]
Type=simple
ExecStart=@PYTHON@ -m deepiri_git_handshake.daemon --foreground
Restart=on-failure
RestartSec=3
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=default.target
