# packaging/launchd/com.deepiri.git-handshake.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.deepiri.git-handshake</string>
    <key>ProgramArguments</key>
    <array>
        <string>@PYTHON@</string>
        <string>-m</string>
        <string>deepiri_git_handshake.daemon</string>
        <string>--foreground</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>@HOME@/Library/Logs/deepiri-git-handshake.log</string>
    <key>StandardErrorPath</key>
    <string>@HOME@/Library/Logs/deepiri-git-handshake.log</string>
</dict>
</plist>
